# deploys/hosts/swarm-mng/files/etc/cron.d/backup.example
00 01  *  *  * root echo "Run swarm backup"
